# cron/cleanup-events.cron
# ToolRadar - EventLog自動クリーンアップ
# 毎週日曜日の深夜3時に30日より古いログを削除

# PATH設定
PATH=/usr/local/bin:/usr/bin:/bin

# MAILTO設定（エラー時の通知先）
MAILTO=admin@toolradar.jp

# 毎週日曜日 3:00 AM JST
0 3 * * 0 cd /home/kwwit/Desktop/trading-tools-platform && docker compose exec -T backend python manage.py cleanup_events >> /var/log/toolradar/cleanup.log 2>&1

# ログローテーション実行（月初1日 4:00 AM）
0 4 1 * * /usr/sbin/logrotate /etc/logrotate.d/toolradar -f
